# crontab
# VATcomply scheduled tasks
# This file should be installed in the container or server running the application

# Load exchange rates every hour at minute 10
10 * * * * . /app/.env.docker; cd /app && python manage.py load_rates --last-90-days >> /proc/1/fd/1 2>&1

# Load countries once daily at 2am
0 2 * * * . /app/.env.docker; cd /app && python manage.py load_countries >> /proc/1/fd/1 2>&1

# Load VAT rates once daily at 3am
0 3 * * * . /app/.env.docker; cd /app && python manage.py load_vat_rates >> /proc/1/fd/1 2>&1
